# monitor_bilge.service
[Unit]
Description=Monitor Bilge

[Service]
Type=simple
ExecStart=/usr/local/sbin/monitor_bilge.py

[Install]
WantedBy=multi-user.target
